# debian.new/postinst
#!/bin/sh

set -e

update-rc.d bwbar defaults >/dev/null
invoke-rc.d bwbar start

#DEBHELPER#

exit 0
